# source/payload/com.user.keepzebawake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.keepzebawake</string>
    <key>ProgramArguments</key>
    <array>
        <string>$HOME/keep_zeb_awake.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/zeb_keepalive.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/zeb_keepalive.log</string>
</dict>
</plist>

